# wmbusmeters-ha-addon-edge/rootfs/etc/services.d/socat/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start socat service
# ==============================================================================

if [ "$(bashio::config 'MbusTCPenabled')" = "yes" ]
then
  bashio::log.info "Starting socat service."
  exec /socat.sh
fi
